# var/lib/lernstick-exam-client/persistent/etc/systemd/system/launch.timer
[Unit]
Description=launch timer

[Timer]
OnUnitInactiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
