# var/lib/lernstick-exam-client/persistent/etc/systemd/system/launch.service
[Unit]
Description=launch

[Service]
Type=oneshot
ExecStart=/usr/bin/launch
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=launch
