# var/lib/lernstick-exam-client/persistent/etc/systemd/system/live_overview.service
[Unit]
Description=live_overview

[Service]
Type=simple
ExecStart=/usr/bin/live_overview

[Install]
WantedBy=graphical.target
